# clonemac
root@R6100:/etc/config# vi /etc/init.d/clonemac

#!/bin/sh /etc/rc.common
# Copyright (C) 2014 OpenWrt.org

START=94
STOP=15

start() {
        ifconfig eth1 down
        ifconfig eth1 hw ether  DC:EF:09:F2:13:70
        ifconfig eth1 up
}           

stop() {    
        echo "Stop."
}      


config 'interface' 'wan'
        option 'ifname' 'eth0.1'
        option 'proto' 'dhcp'
        option 'macaddr' 'd8:64:c7:5e:f1:c6'
        
or rc.local
